# samba/root/etc/services.d/samba/run
#!/usr/bin/with-contenv bash

smbd --foreground --log-stdout --configfile=/config/smb.conf