# nginx/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

nginx -V
nginx -g "daemon off; error_log /dev/stdout info;" -c /config/nginx.conf