# avahi/root/etc/services.d/avahi/run
#!/usr/bin/with-contenv bash

avahi-daemon --debug -f /config/avahi-daemon.conf
